# BuildTools/jenkins-slave-dotnetcore/supervisord.conf
[supervisord]
nodaemon=true

[program:dockerd-entrypoint]
command=dockerd-entrypoint.sh

[program:jenkins-slave]
command=jenkins-slave
